# station3/docs/sysConfig/bird-startup.service
# put in /etc/systemd/system/, 'systemctl enable bird-startup.service', waits for network up, in contrast to crontab or rc.local (latter depends on systemctl enable rc-local.service)
# teste zuvor: which bash, which sleep, which curl usw.
# copy the 'Environment=PATH=' from 'echo $PATH' in login shell, where it is set in '.bashrc'. In both the birdvenv should be included.
# REMOVE ALL COMMENTS BEFORE INSTALLING
[Unit]
Description=Birds Startup Script
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/station3

# Allow subprocesses to stay alive after the main script ends
# KillMode=process
# Optional alternative: KillMode=none if you want *nothing* killed
KillMode=none

ExecStartPre=/usr/bin/sleep 5
ExecStart=/usr/bin/bash /home/pi/station3/startup.sh

StandardOutput=append:/home/pi/station3/logs/startup.service.log
StandardError=append:/home/pi/station3/logs/startup.service.err

Environment=PATH=/home/pi/birdvenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VIRTUAL_ENV=/home/pi/birdvenv
Environment=HOME=/home/pi
Environment=PYTHONUNBUFFERED=1

Restart=on-failure
RestartSec=5s

# Don’t timeout early
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
