# station3/docs/sysConfig/bird-startup2.service
[Unit]
Description=Bird Station Startup Supervisor
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/station3

# Run a small supervision loop so systemd can restart us
ExecStart=/usr/bin/bash -lc '
  while true; do
    echo "$(date): Starting bird scripts" >> /home/pi/station3/logs/supervisor.log 2>&1
    setsid /home/pi/activate_birdvenv.sh && setsid /home/pi/station3/startup.sh >> /home/pi/station3/logs/startup.service.log 2>&1
    echo "$(date): bird scripts exited, restarting in 10s..." >> /home/pi/station3/logs/supervisor.log 2>&1
    sleep 10
  done
'

# Restart if the supervision loop dies (very unlikely)
Restart=always
RestartSec=5

# Faster boot
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
